# net/nylon/files/nylon.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2006-2012 OpenWrt.org

START=50
SERVICE_USE_PID=1

NAME="nylon"
DAEMON="/usr/bin/${NAME}"
CONFIG="/etc/${NAME}.conf"

start() {
	if [ -r ${CONFIG} ]
	then
		service_start ${DAEMON} -c ${CONFIG}
	fi
}

stop() {
	service_stop ${DAEMON}
}
